# services/iotsentinel-provision.service
[Unit]
Description=IoTSentinel WiFi Provisioning
After=NetworkManager.service
Before=iotsentinel-dashboard.service
Wants=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Pick up IOTSENTINEL_WIFI_COUNTRY (set by the wizard) so the setup hotspot uses the
# user's Wi-Fi region. Optional ('-'): absent on first boot, present after setup.
EnvironmentFile=-/home/sentinel/iotsentinel/.env
# Generous: setup_hotspot.sh waits for NetworkManager to take over wlan0 on a cold
# boot and then retries the AP bring-up a few times, so the worst case is ~2 min.
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal
ExecStart=/bin/bash -c '\
  APP="/home/sentinel/iotsentinel"; ENVF="$APP/.env"; \
  if ! grep -q "^FLASK_SECRET_KEY=" "$ENVF" 2>/dev/null; then \
    echo "FLASK_SECRET_KEY=$(openssl rand -hex 32)" >> "$ENVF"; \
    chown "$(stat -c %U:%G "$APP" 2>/dev/null || echo sentinel:sentinel)" "$ENVF" 2>/dev/null || true; \
    echo "[provision] Generated a per-device FLASK_SECRET_KEY"; \
  fi; \
  bash "$APP/scripts/setup_hotspot.sh" boot'

[Install]
WantedBy=multi-user.target
