# services/iotsentinel-dashboard.service
[Unit]
Description=IoTSentinel Web Dashboard
After=iotsentinel-backend.service
BindsTo=iotsentinel-backend.service

[Service]
Type=simple
User=sentinel
WorkingDirectory=/home/sentinel/iotsentinel
ExecStart=/home/sentinel/iotsentinel/venv/bin/python3 /home/sentinel/iotsentinel/dashboard/app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
