# services/iotsentinel-zeek-parser.service
[Unit]
Description=IoTSentinel Zeek Log Parser
After=network.target
Requires=network.target

[Service]
Type=simple
User=sentinel
WorkingDirectory=/home/sentinel/iotsentinel
ExecStart=/home/sentinel/iotsentinel/venv/bin/python3 /home/sentinel/iotsentinel/capture/zeek_log_parser.py --watch --interval 60
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iotsentinel-zeek

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target