# services/iotsentinel-backend.service
[Unit]
Description=IoTSentinel Backend Orchestrator (Parser and ML Engine)
# Wait for NetworkManager so the home-Wi-Fi uplink and (gateway mode) the AP
# interface exist before the orchestrator brings the access point up and points
# Zeek at it.
After=network-online.target NetworkManager.service
Wants=network-online.target
# Never permanently give up. A transient disk hiccup (e.g. a full/recovering SD
# card making the DB connect fail) would otherwise hit systemd's DEFAULT start
# limit (5 restarts / 10s) and leave the service dead until a manual reboot.
# With the limit disabled the orchestrator keeps retrying and self-heals once
# the filesystem recovers (or fsck.repair runs at the next boot).
StartLimitIntervalSec=0

[Service]
Type=simple
User=sentinel
WorkingDirectory=/home/sentinel/iotsentinel
ExecStart=/home/sentinel/iotsentinel/venv/bin/python3 /home/sentinel/iotsentinel/orchestrator.py
Restart=on-failure
# 10s (was 5s) so a crash-loop does not spam the journal and fill the disk.
RestartSec=10
LimitNOFILE=8192
# NOTE: deliberately no MemoryMax -- the Pi kernel cmdline sets
# cgroup_disable=memory, so memory cgroup limits are a no-op on this hardware.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
